# crontab
0 */4 * * * /bin/bash /app/workflow.sh >> /app/cron.log 2>&1
